# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start
autorestart=false
